# crontab/crontab
# backup todo list every hour
PATH=~/bin:/usr/bin:/bin:~/.local/bin
0 * * * * backup-todo-list
